# src/config/systemd/user/alacritty@.service
[Unit]
Description=alacritty - Terminal emulator
Documentation=man:alacritty(5)

[Service]
Type=simple
ExecStart=alacritty
