# src/config/systemd/user/firefox@.service
[Unit]
Description=Firefox - Web Browser
After=graphical-session.target
Requires=graphical-session.target

[Service]
ExecStart=firefox
