# src/config/systemd/user/krita@.service
[Unit]
Description=Krita drawing application
After=graphical-session.target
Requires=graphical-session.target

[Service]
ExecStart=krita
