# src/config/systemd/user/mypaint@.service
[Unit]
Description=MyPaint drawing application
After=graphical-session.target
Requires=graphical-session.target

[Service]
ExecStart=mypaint
